# root/etc/services.d/aria2/run
#! /usr/bin/with-contenv bash

#è®¾ç½®umask
umask "$UMASK"

#å¯åŠ¨aria2
exec s6-setuidgid aria2 aria2c --conf-path=/config/aria2.conf > /dev/stdout 2 > /dev/stderr
